# packaging/powerstats-user.service
[Unit]
Description=PowerStats Telemetry Daemon
Documentation=https://github.com/powerstats/powerstats
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/share/powerstats/daemon.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=powerstats-daemon

[Install]
WantedBy=default.target
