# packaging/powerstats.service
[Unit]
Description=PowerStats Telemetry Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/share/powerstats/daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
